# lib/redis/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/opt/redis-3.0.0-rc1/src/redis-server /opt/redis.conf